# pi-cluster/cloudflare/cloudflare_dns.service
[Unit]
Description=Cloudflare dynamic IP updating

[Service]
WorkingDirectory=/home/pi/cloudflare
User=pi
ExecStart=/home/pi/cloudflare/checkip.sh

